# agents/systemd/claude-watchdog-itachi.service
[Unit]
Description=Watchdog for Itachi (Claude Agent) — Konoha SSE → tmux/terminal
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Environment=HOME=/home/YOUR_USERNAME
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/home/YOUR_USERNAME/.local/bin
# For WSL/remote agents: use domain (direct IP not accessible from outside)
Environment=KONOHA_URL=https://agent.eaprelsky.ru
# Set your agent token (obtained after POST /agents/register with invite token)
Environment=KONOHA_TOKEN=YOUR_AGENT_TOKEN
Environment=ITACHI_AGENT=itachi
Environment=ITACHI_TMUX=itachi
ExecStart=/usr/bin/python3 /home/YOUR_USERNAME/scripts/watchdog-itachi.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
