# agents/systemd/claude-mirai.service
[Unit]
Description=Claude Mirai (Agent #3 - Email/Data Processor)
After=network.target redis-server.service

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/scripts/claude-mirai-service.sh
Restart=always
RestartSec=10
Environment=HOME=/home/ubuntu
Environment=KONOHA_AGENT_ID=mirai
Environment=PATH=/home/ubuntu/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.local/bin

[Install]
WantedBy=multi-user.target
