# systemd/claude-kiba.service
[Unit]
Description=Claude Agent Kiba (System Guardian)
After=network.target redis.service
Wants=claude-watchdog-kiba.service
After=claude-watchdog-kiba.service

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/scripts/claude-kiba-service.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
