# systemd/claude-watchdog-mirai.service
[Unit]
Description=Watchdog for Mirai (Claude Agent #3) — Konoha SSE → tmux
After=network.target konoha.service

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/home/ubuntu/.agent-env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/scripts/watchdog-mirai.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/home/ubuntu/.local/bin

[Install]
WantedBy=multi-user.target
