# systemd/claude-ino.service
[Unit]
Description=Claude Agent Ino (Nocturna Marketing) — on-demand
After=network.target redis.service
Wants=claude-watchdog-ino.service
After=claude-watchdog-ino.service

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/scripts/claude-ino-service.sh
ExecStop=/usr/bin/tmux kill-session -t ino
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
