# systemd/claude-watchdog-guy.service
[Unit]
Description=Claude Watchdog Guy (Konoha SSE → tmux)
After=network.target redis.service

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/home/ubuntu/.agent-env
ExecStart=/usr/bin/python3 /home/ubuntu/scripts/watchdog-guy.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
