# systemd/claude-watchdog-kakashi.service
[Unit]
Description=Claude Watchdog Kakashi (Konoha SSE + GitHub Issues → tmux)
After=network.target redis.service

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/home/ubuntu/.agent-env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/scripts/watchdog-kakashi.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
